# service-gateway/launchd/com.teamshell.service-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.teamshell.service-gateway</string>

  <!-- 按本机实际路径修改下面三处：node 路径（nvm）、仓库路径、PATH -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhengweibin/.nvm/versions/node/v24.14.0/bin/node</string>
    <string>/Users/zhengweibin/Desktop/team-shell/data-analysis/service-gateway/dist/service.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/zhengweibin/Desktop/team-shell/data-analysis</string>

  <!-- launchd 不继承 shell 环境：PATH 必须显式包含 node（nvm）与 lark-cli -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/zhengweibin/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/zhengweibin/Library/Logs/service-gateway/out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zhengweibin/Library/Logs/service-gateway/err.log</string>
</dict>
</plist>
